# tap/deploy/tap-003/skylens-tap.service
[Unit]
Description=Skylens Tap Drone Detection
After=skylens-nats-proxy.service
Requires=skylens-nats-proxy.service

[Service]
Type=notify
WatchdogSec=30
ExecStart=/opt/skylens/skylens-tap /opt/skylens/config.toml /opt/skylens/drone_models.json
Restart=on-failure
RestartSec=10
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
