# tap/deploy/tap-003/tailscaled-skylens.service
[Unit]
Description=Tailscale (Skylens) - userspace networking
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/tailscaled --state=/var/lib/tailscale-skylens/tailscaled.state --socket=/var/run/tailscale-skylens/tailscaled.sock --tun=userspace-networking --port=0
RuntimeDirectory=tailscale-skylens
StateDirectory=tailscale-skylens
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
